# linux/.config/autostart/dunst.desktop
[Desktop Entry]
Encoding=UTF-8
Name=dunst
Comment=notification
Exec=dunst
Type=Application
