# linux/.config/autostart/mpd.desktop
[Desktop Entry]
Encoding=UTF-8
Name=mpd
Comment=music
Exec=mpd
Type=Application
