# linux/.config/autostart/sxhkd.desktop
[Desktop Entry]
Encoding=UTF-8
Name=sxhkd
Comment=hotkey
Exec=sxhkd
Type=Application
